# cloudflare-ddns/cloudflare-ddns.service
[Unit]
Description=Cloudflare DDNS Updater

[Service]
Type=oneshot
ExecStart=/bin/bash __PATH__/cloudflare-ddns/update.sh
